# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/home/ubuntu/analysis/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=5
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
user=ubuntu

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ── Agent servers ──────────────────────────────────────────────────────────────

[program:agent-profiler]
command=/home/ubuntu/analysis/venv/bin/python agent_servers/server_base.py --agent profiler --port 8001
directory=/home/ubuntu/analysis
user=ubuntu
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/home/ubuntu/analysis/logs/profiler.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/home/ubuntu/analysis/logs/profiler.log
stderr_logfile_maxbytes=5MB
environment=USE_A2A_MULTISERVER="true",PYTHONIOENCODING="utf-8",PYTHONUTF8="1"

[program:agent-discovery]
command=/home/ubuntu/analysis/venv/bin/python agent_servers/server_base.py --agent discovery --port 8002
directory=/home/ubuntu/analysis
user=ubuntu
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/home/ubuntu/analysis/logs/discovery.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/home/ubuntu/analysis/logs/discovery.log
stderr_logfile_maxbytes=5MB
environment=USE_A2A_MULTISERVER="true",PYTHONIOENCODING="utf-8",PYTHONUTF8="1"

[program:agent-coder]
command=/home/ubuntu/analysis/venv/bin/python agent_servers/server_base.py --agent coder --port 8003
directory=/home/ubuntu/analysis
user=ubuntu
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/home/ubuntu/analysis/logs/coder.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/home/ubuntu/analysis/logs/coder.log
stderr_logfile_maxbytes=5MB
environment=USE_A2A_MULTISERVER="true",PYTHONIOENCODING="utf-8",PYTHONUTF8="1"

[program:agent-synthesis]
command=/home/ubuntu/analysis/venv/bin/python agent_servers/server_base.py --agent synthesis --port 8004
directory=/home/ubuntu/analysis
user=ubuntu
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/home/ubuntu/analysis/logs/synthesis.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/home/ubuntu/analysis/logs/synthesis.log
stderr_logfile_maxbytes=5MB
environment=USE_A2A_MULTISERVER="true",PYTHONIOENCODING="utf-8",PYTHONUTF8="1"

[program:agent-critic]
command=/home/ubuntu/analysis/venv/bin/python agent_servers/server_base.py --agent critic --port 8005
directory=/home/ubuntu/analysis
user=ubuntu
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/home/ubuntu/analysis/logs/critic.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/home/ubuntu/analysis/logs/critic.log
stderr_logfile_maxbytes=5MB
environment=USE_A2A_MULTISERVER="true",PYTHONIOENCODING="utf-8",PYTHONUTF8="1"

[program:agent-dag-builder]
command=/home/ubuntu/analysis/venv/bin/python agent_servers/server_base.py --agent dag_builder --port 8006
directory=/home/ubuntu/analysis
user=ubuntu
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/home/ubuntu/analysis/logs/dag_builder.log
stdout_logfile_maxbytes=5MB
stderr_logfile=/home/ubuntu/analysis/logs/dag_builder.log
stderr_logfile_maxbytes=5MB
environment=USE_A2A_MULTISERVER="true",PYTHONIOENCODING="utf-8",PYTHONUTF8="1"

; ── Main FastAPI backend (port 8000) ──────────────────────────────────────────

[program:adk]
command=/home/ubuntu/analysis/venv/bin/python -X utf8 main.py
directory=/home/ubuntu/analysis
user=ubuntu
autostart=true
autorestart=true
startretries=5
startsecs=15
stdout_logfile=/home/ubuntu/analysis/logs/main.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/home/ubuntu/analysis/logs/main.log
stderr_logfile_maxbytes=10MB
environment=USE_A2A_MULTISERVER="true",PYTHONIOENCODING="utf-8",PYTHONUTF8="1"

; ── Group all processes so one command restarts everything ─────────────────────

[group:all-agents]
programs=agent-profiler,agent-discovery,agent-coder,agent-synthesis,agent-critic,agent-dag-builder

[group:full-stack]
programs=agent-profiler,agent-discovery,agent-coder,agent-synthesis,agent-critic,agent-dag-builder,adk
